# training/systemd/vacation-overflow-watch@.timer
[Unit]
Description=Check whether Blood Bowl vacation overflow %i may start

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true
Unit=vacation-overflow-watch@%i.service

[Install]
WantedBy=timers.target
